# Oldest/systemd/manage-swap.service
[Unit]
Description=Manage Swap Files Dynamically
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/manage_swap.sh start
ExecStop=/usr/local/bin/manage_swap.sh stop
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
